# My_challenges/Custom/BurfieSuite/supervisord.conf
[supervisord]
nodaemon=true

[program:flask]
command=python3 ape.py
autostart=true
autorestart=true
user=root
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stdout

[program:vsftpd]
command=vsftpd